# smmonitor-agent/lib/systemd/system/patrold.service
[Unit]
Description=Patrol Monitoring Agent
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/patrold
Restart=always
RestartSec=5
StartLimitBurst=10
StartLimitInterval=60

[Install]
WantedBy=multi-user.target