# smmonitor-agent/DEBIAN/postinst
#!/bin/bash
set -e

# Ask for API key and URL if not already configured
echo -n "Enter your Patrol Agent API key: "
read API_KEY

echo -n "Enter the Service Monitor URL (e.g., https://your-monitor.example.com): "
read API_BASE_URL

# Update config with user input
CONFIG_FILE="/etc/patrold/config.json"
TMP_CONFIG=$(mktemp)
cat $CONFIG_FILE | jq ".api_key = \"$API_KEY\" | .api_url = \"$API_BASE_URL\"" > $TMP_CONFIG
cat $TMP_CONFIG > $CONFIG_FILE
rm $TMP_CONFIG

# Fix permissions
chmod 644 $CONFIG_FILE

# Enable and start the service
systemctl daemon-reload
systemctl enable patrold.service
systemctl start patrold.service

# Print success message
echo "=============================================="
echo "Patrol Agent has been installed and started."
echo "Configuration stored in: $CONFIG_FILE"
echo "To check the status: systemctl status patrold.service"
echo "To view logs: journalctl -u patrold.service -f"
echo "=============================================="

# Exit with success
exit 0