# smmonitor-agent/DEBIAN/postinst
#!/bin/bash
set -e

# Enable and start the service
systemctl daemon-reload
systemctl enable patrold.service
systemctl start patrold.service

# Print success message
echo "Patrol Agent has been installed and started."
echo "To check the status, run: systemctl status patrold.service"
echo "To view logs, run: journalctl -u patrold.service"

# Exit with success
exit 0