# smmonitor-agent/DEBIAN/prerm
#!/bin/bash
set -e

# Stop and disable the service
systemctl stop patrold.service || true
systemctl disable patrold.service || true

# Exit with success
exit 0